# service/log/run
#!/bin/sh
exec setuidgid sugyan multilog t n10 .
